# bot/systemd/otel-collector.service
[Unit]
Description=OpenTelemetry Collector (Monium)
After=network.target
Before=cards-order-bot.service

[Service]
Type=simple
User=mbabaev
WorkingDirectory=/home/mbabaev/CardsOrder
EnvironmentFile=/home/mbabaev/CardsOrder/bot/.env
ExecStart=/usr/local/bin/otelcol --config /home/mbabaev/CardsOrder/otel-collector.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
